# insight_calc/pluto-insight.service
[Unit]
Description=Pluto Remote Server for InsightCircle
After=network.target

[Service]
Type=simple
User=gcr
WorkingDirectory=/home/gcr/populi.Wk/InsightCircle/insight_calc
# We point directly to the Julia binary and pass the project flag
ExecStart=/home/gcr/.juliaup/bin/julia --project=@. /home/gcr/populi.Wk/InsightCircle/insight_calc/pluto-remote.jl
Restart=on-failure
# Load your environment variables here
Environment="INSIGHT_PLUTO_PORT=5202"
Environment="INSIGHT_PLUTO_SECRET=your_secure_string_here"

[Install]
WantedBy=multi-user.target